# profile.d/clojure.sh
# clojure
if [ -e /usr/local/bin/clj ]; then
    export CLASSPATH=${CLASSPATH}:${HOME}/.lein/plugins
fi


