# profile.d/lfe.sh
LFE_VER=0.8
. ${HOME}/activate

export PATH=$PATH:~/lib/lfe-${LFE_VER}/bin
